# meshtastic-dashboard.service
# Example unit for a git checkout cloned to /opt/meshyface.
# scripts/deploy_meshyface.sh renders a host-specific unit during --bootstrap.
[Unit]
Description=Meshyface Meshtastic Web Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=meshyface
Group=dialout
WorkingDirectory=/opt/meshyface
EnvironmentFile=/etc/meshyface/dashboard.env
ExecStart=/opt/meshyface/.venv/bin/python /opt/meshyface/mesh_dashboard.py --http-host 0.0.0.0 --http-port 8877 --refresh-ms 3000
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
